# docker/supervisor/reverb.conf
[supervisord]
nodaemon=true

[program:laravel-app]
command=php artisan serve --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=www-data
directory=/var/www
stdout_logfile=/var/log/supervisor/app.log
stderr_logfile=/var/log/supervisor/app-error.log

[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=8000
autostart=true
autorestart=true
user=www-data
directory=/var/www
stdout_logfile=/var/log/supervisor/reverb.log
stderr_logfile=/var/log/supervisor/reverb-error.log
